# deploy/aipplify-jobs-sync.service
[Unit]
Description=Aipplify — синхронизация вакансий (парсеры в БД)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/Aipplify/Front

# Подставьте свой путь к проекту при установке на другой сервер.
EnvironmentFile=-/root/Aipplify/Front/.env

# Опционально (перекрывают значения по умолчанию в скрипте):
# Environment=JOBS_SYNC_DEGEN_MAX_PAGES=15
# Environment=JOBS_SYNC_SKIP_CRYPTO=1

# Используем bash -lc, чтобы подтянуть PATH (nvm/fnm), если npm не в /usr/bin
ExecStart=/bin/bash -lc 'cd /root/Aipplify/Front && npm run jobs:sync'

# Логи
StandardOutput=append:/var/log/aipplify-jobs-sync.log
StandardError=append:/var/log/aipplify-jobs-sync.log

User=root

[Install]
WantedBy=multi-user.target
